# configs/etc/supervisord.d/cobblerd.conf
[program:cobblerd]
command=/usr/local/bin/cobblerd -F
priority=3
startsecs=10
autorestart=true
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0